# automations/launchd/com.needthisdone.auto-eval.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.needthisdone.auto-eval</string>

    <key>ProgramArguments</key>
    <array>
        <string>/opt/homebrew/bin/npx</string>
        <string>tsx</string>
        <string>/Users/abereyes/Projects/Personal/Need_This_Done/automations/run-eval.ts</string>
    </array>

    <key>WorkingDirectory</key>
    <string>/Users/abereyes/Projects/Personal/Need_This_Done</string>

    <!-- Run every 4 hours (14400 seconds) -->
    <key>StartInterval</key>
    <integer>14400</integer>

    <!-- Environment variables for Claude CLI and Node -->
    <key>EnvironmentVariables</key>
    <dict>
        <key>PATH</key>
        <string>/Users/abereyes/.local/bin:/opt/homebrew/bin:/usr/local/bin:/usr/bin:/bin</string>
        <key>HOME</key>
        <string>/Users/abereyes</string>
    </dict>

    <!-- Log output -->
    <key>StandardOutPath</key>
    <string>/Users/abereyes/Projects/Personal/Need_This_Done/automations/launchd/auto-eval.stdout.log</string>
    <key>StandardErrorPath</key>
    <string>/Users/abereyes/Projects/Personal/Need_This_Done/automations/launchd/auto-eval.stderr.log</string>

    <!-- Don't restart on failure -->
    <key>KeepAlive</key>
    <false/>

    <!-- Run even if user is not logged in (for Mac Mini server) -->
    <key>RunAtLoad</key>
    <true/>
</dict>
</plist>
